# files/kodi.supervisor.conf
[program:kodi]
user         = force
environment  = DISPLAY=":0",HOME="/home/force"
command      = /usr/bin/kodi
autostart    = true
autorestart  = true

; vim: set filetype=dosini:
